# dmonitor.service
[Unit]
Description=dmonitor

[Service]
User=dmonitor
ExecStart=/etc/init.d/dmonitor
Restart=always

[Install]
WantedBy=multi-user.target
